# minus-update.timer
[Unit]
Description=Daily Minus update at 6am ET

[Timer]
# 6am Eastern Time (handles DST automatically)
OnCalendar=*-*-* 06:00:00 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
